# templates/zrepl-kopia-uploader.service.j2
[Service]
Type=oneshot
ExecStart=/usr/local/bin/zrepl-kopia-uploader
