# templates/zrepl-kopia-uploader.timer.j2
[Unit]
Description=Backup zrepl snapshots with kopia

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
